# launchagent/com.noelmom.vox-helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.noelmom.vox-helper</string>

  <key>AssociatedBundleIdentifiers</key>
  <array>
    <string>com.noelmom.vox-helper</string>
  </array>

  <!-- TEMPORARY: runs python3 directly from the permanent venv until VoxHelper.app
       is code-signed with a Developer ID certificate. Once signed, revert to the
       .app bundle approach — see BACKLOG.md "Revert helper to signed VoxHelper.app". -->
  <key>ProgramArguments</key>
  <array>
    <string>VOX_APP_SUPPORT/venv/bin/python3</string>
    <string>VOX_APP_SUPPORT/menubar/vox_helper.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>VOX_APP_SUPPORT</string>

  <key>StandardOutPath</key>
  <string>VOX_LOG_DIR/vox-helper.log</string>
  <key>StandardErrorPath</key>
  <string>VOX_LOG_DIR/vox-helper-error.log</string>

  <!--
    RunAtLoad = true — the menu bar helper starts automatically on login.
    This is intentional: the helper is the UI layer, not the server.
    The server (com.noelmom.vox) remains manual-start.
  -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Restart if the helper crashes -->
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <!-- Required for rumps to attach to the macOS window server -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

</dict>
</plist>
