# launchagent/com.melolabdev.vox-helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.melolabdev.vox-helper</string>

  <!-- Executable lives inside ~/Applications/VoxHelper.app — no admin
       privileges required. The helper script is copied to Application
       Support on install so it survives the project folder being deleted. -->
  <key>ProgramArguments</key>
  <array>
    <string>VOX_HOME/Applications/VoxHelper.app/Contents/MacOS/vox-helper</string>
    <string>VOX_APP_SUPPORT/menubar/vox_helper.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>VOX_APP_SUPPORT</string>

  <key>StandardOutPath</key>
  <string>VOX_LOG_DIR/vox-helper.log</string>
  <key>StandardErrorPath</key>
  <string>VOX_LOG_DIR/vox-helper-error.log</string>

  <!--
    RunAtLoad = true — the menu bar helper starts automatically on login.
    This is intentional: the helper is the UI layer, not the server.
    The server (com.melolabdev.vox) remains manual-start.
  -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Restart if the helper crashes -->
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <!-- Required for rumps to attach to the macOS window server -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

</dict>
</plist>
