# launchagent/com.melolabdev.vox.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.melolabdev.vox</string>

  <!-- run.sh is written into Application Support by install-agent.sh.
       The server runs entirely from APP_SUPPORT — no project folder dependency. -->
  <key>ProgramArguments</key>
  <array>
    <string>VOX_APP_SUPPORT/VoxServer.app/Contents/MacOS/vox-server</string>
  </array>

  <key>WorkingDirectory</key>
  <string>VOX_APP_SUPPORT</string>

  <key>StandardOutPath</key>
  <string>VOX_LOG_DIR/vox.log</string>
  <key>StandardErrorPath</key>
  <string>VOX_LOG_DIR/vox-error.log</string>

  <!--
    RunAtLoad = false — server does NOT start automatically on login.
    SHIPPING NOTE: change to <true/> when packaging the .app.
  -->
  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>TimeOut</key>
  <integer>120</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>VOX_HOST</key>
    <string>0.0.0.0</string>
    <key>VOX_PORT</key>
    <string>8000</string>
    <key>VOX_DEVICE</key>
    <string>auto</string>
  </dict>

</dict>
</plist>
